# bin/polyglot.service
[Unit]
Description=Polyglot service allows for converting between a variety of file types

[Service]
Type=simple
WorkingDirectory=/home/polyglot/polyglot
ExecStart=/usr/bin/java/ -cp polyglot.jar:lib/* -Xmx1g edu.illinois.nca.isda.softwareserver.polyglot.PolyglotRestlet
ExecStop=/bin/kill -INT $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
